# rootfs/etc/inittab
::sysinit:/etc/init.d/rcS
ttyAMA0::askfirst:/bin/ash
